# .crontab
MAILTO=dbarkats@cfa.harvard.edu
00 1,2,3,5,6,7,9,10,11,13,14,15,17,18,19,21,22,23 * * *  nice -n -19 $HOME/WVR_Omnisys/wvrObserve1hr.py &>>  $HOME/WVR_Omnisys/data_tmp/wvrLog.txt
00 0,4,8,12,16,20 * * * nice -n -19 $HOME/WVR_Omnisys/wvrNoise.py &>>  $HOME/WVR_Omnisys/data_tmp/wvrLog.txt
10 * * * *   /home/dbarkats/WVR_Omnisys/wvrArchiveData.py 
00 0 * * * $HOME/WVR_Omnisys/KeepSerPortAlive.py -f &>> $HOME/WVR_Omnisys/data_tmp/wvrLog.txt
